# src/etc/systemd/system/serial-getty@.service.d/50-oem-autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --keep-baud --autologin oem --noclear 1500000,115200,57600,38400,9600 %I $TERM
